# deploy/systemd/dayboard-postgres-backup.timer
[Unit]
Description=Run the Dayboard PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=15m
Unit=dayboard-postgres-backup.service

[Install]
WantedBy=timers.target
