# deploy/systemd/dayboard-web.service
[Unit]
Description=Dayboard Next.js web service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/dayboard/apps/web
Environment=NODE_ENV=production
Environment=NEXT_PUBLIC_DAYBOARD_API_BASE_URL=/dayboard-api
Environment=NEXT_PUBLIC_DAYBOARD_BASE_PATH=/dayboard
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3001
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
