# deploy/systemd/dayboard-api.service
[Unit]
Description=Dayboard FastAPI service
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/dayboard/apps/api
EnvironmentFile=/root/dayboard/.env
ExecStart=/usr/local/bin/uv run fastapi run src/dayboard/main.py --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
