# deploy/systemd/dayboard-worker.service
[Unit]
Description=Dayboard arq command worker
After=network-online.target docker.service dayboard-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/dayboard/apps/api
EnvironmentFile=/root/dayboard/.env
ExecStart=/usr/local/bin/uv run arq dayboard.workers.commands.WorkerSettings
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
